# tablas/archivos_tabla/osquery_schedule2.json

  average_memory:0:blacklisted:0:executions:0:interval:300:last_executed:0:name:crontab:output_size:0:query:SELECT * FROM crontab;:system_time:0:user_time:0:wall_time:0:
  average_memory:0:blacklisted:0:executions:0:interval:3600:last_executed:0:name:system_info:output_size:0:query:SELECT hostname: cpu_brand: physical_memory FROM system_info;:system_time:0:user_time:0:wall_time:0:
  average_memory:0:blacklisted:0:executions:0:interval:3600:last_executed:0:name:system_profile:output_size:0:query:SELECT * FROM osquery_schedule;:system_time:0:user_time:0:wall_time:0:
  average_memory:0:blacklisted:0:executions:0:interval:300:last_executed:0:name:pack_fim_file_events:output_size:0:query:select * from file_events;:system_time:0:user_time:0:wall_time:0:
  average_memory:0:blacklisted:0:executions:0:interval:86400:last_executed:0:name:pack_osquery-monitoring_events:output_size:0:query:select name: publisher: type: subscriptions: events: active from osquery_events;:system_time:0:user_time:0:wall_time:0:
  average_memory:0:blacklisted:0:executions:0:interval:600:last_executed:0:name:pack_osquery-monitoring_osquery_info:output_size:0:query:select i.*: p.resident_size: p.user_time: p.system_time: time.minutes as counter from osquery_info i: processes p: time where p.pid = i.pid;:system_time:0:user_time:0:wall_time:0:
  average_memory:0:blacklisted:0:executions:0:interval:7200:last_executed:0:name:pack_osquery-monitoring_schedule:output_size:0:query:select name: interval: executions: output_size: wall_time: (user_time\/executions) as avg_user_time: (system_time\/executions) as avg_system_time: average_memory: last_executed from osquery_schedule;:system_time:0:user_time:0:wall_time:0:
  average_memory:0:blacklisted:0:executions:0:interval:3600:last_executed:0:name:pack_incident-response_arp_cache:output_size:0:query:select * from arp_cache;:system_time:0:user_time:0:wall_time:0:
  average_memory:0:blacklisted:0:executions:0:interval:3600:last_executed:0:name:pack_incident-response_crontab:output_size:0:query:select * from crontab;:system_time:0:user_time:0:wall_time:0:
  average_memory:0:blacklisted:0:executions:0:interval:86400:last_executed:0:name:pack_incident-response_disk_encryption:output_size:0:query:select * from disk_encryption;:system_time:0:user_time:0:wall_time:0:
  average_memory:0:blacklisted:0:executions:0:interval:86400:last_executed:0:name:pack_incident-response_etc_hosts:output_size:0:query:select * from etc_hosts;:system_time:0:user_time:0:wall_time:0:
  average_memory:0:blacklisted:0:executions:0:interval:3600:last_executed:0:name:pack_incident-response_ip_forwarding:output_size:0:query:select * from system_controls where oid = '4.30.41.1' union select * from system_controls where oid = '4.2.0.1';:system_time:0:user_time:0:wall_time:0:
  average_memory:0:blacklisted:0:executions:0:interval:3600:last_executed:0:name:pack_incident-response_iptables:output_size:0:query:select * from iptables;:system_time:0:user_time:0:wall_time:0:
  average_memory:0:blacklisted:0:executions:0:interval:3600:last_executed:0:name:pack_incident-response_kernel_modules:output_size:0:query:select * from kernel_modules;:system_time:0:user_time:0:wall_time:0:
  average_memory:0:blacklisted:0:executions:0:interval:3600:last_executed:0:name:pack_incident-response_last:output_size:0:query:select * from last;:system_time:0:user_time:0:wall_time:0

